# buildroot/board/raspberrypi0/ovl/etc/profile.d/custom.sh
# Access to personal scripts and business code
PATH=$PATH:/usr/local/bin

# Prompt indicating host name and current directory.
PS1='\h[\W]\$ '
